# src/cron/daily.cron
# m h  dom mon dow   command
10 01 * * * /usr/local/cron_exec.sh