# src/cron/weekly.cron
# m h  dom mon dow   command
10 01 * * 6 /usr/local/cron_exec.sh