# src/cron/weekly.cron
# m h  dom mon dow   command
30 18 * * 2 /usr/local/cron_exec.sh